# rabbitmq_couchdb_email/mqcouchemail.conf
start on started rabbitmq-server
exec sh /home/rt7/mqcouchemail/start_MQCouchEmailSender.sh

